# service/shell-generator-server.service
[Unit]
Description=Shell Generator Server for MCP Center
After=network.target

[Service]
User=root
WorkingDirectory=/usr/lib/euler-copilot-framework
ExecStart=/usr/bin/python3 mcp_center/servers/shell_generator/src/server.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
KillMode=control-group   
KillSignal=SIGTERM      
TimeoutStopSec=30    
SuccessExitStatus=143 

[Install]
WantedBy=multi-user.target
    